# caddy-pro/DEBIAN/postrm
#!/bin/bash

sudo systemctl disable caddyd
sudo systemctl stop caddyd
rm -f '/usr/bin/caddy-pro'
rm -f '/usr/bin/caddy'