# Frpc/DEBIAN/postrm
#!/bin/bash


sudo systemctl disable frpcd
sudo systemctl stop frpcd
rm -f '/usr/bin/frpc'