# frps-pro/DEBIAN/postrm
#!/bin/bash

sudo systemctl disable frpsd
sudo systemctl stop frpsd
rm -f '/usr/bin/frps-pro'
rm -f '/usr/bin/frps'